# deploy/systemd/intrusion-backend.service
[Unit]
Description=Intrusion Detection Backend
After=network.target

[Service]
Type=simple
User=pi
WorkingDirectory=/home/pi/intrusion-suite/backend
Environment="PATH=/home/pi/intrusion-suite/backend/.venv/bin"
ExecStart=/home/pi/intrusion-suite/backend/.venv/bin/uvicorn app.main:app --host 0.0.0.0 --port 8000
Restart=always
RestartSec=10

[Install]
WantedBy=multi-user.target

